# service.conf
[program:oa]
command=/home/vicky/oa/env/bin/gunicorn -b localhost:8000 -w 3 server:app
directory=/home/vicky/oa
user=vicky
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
